# bin/systemd/chrismeyers-info-api.service
[Unit]
Description=chrismeyers.info API service

[Service]
User=chris
ExecStart=/usr/bin/dart /home/chris/www/chrismeyers.info/api/bin/main.dart
WorkingDirectory=/home/chris/www/chrismeyers.info/api/bin
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=chrismeyers-info-api

[Install]
WantedBy=multi-user.target
